# com.watchscraper.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.watchscraper.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jameshuang/Projects/watchScraper/.venv/bin/python</string>
        <string>-m</string>
        <string>watchscraper.cli</string>
        <string>schedule</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jameshuang/Projects/watchScraper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jameshuang/Projects/watchScraper/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/jameshuang/Projects/watchScraper/src</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jameshuang/Library/Logs/watchscraper/scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jameshuang/Library/Logs/watchscraper/scheduler.log</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
